# pdns/etc/supervisor/conf.d/pdns_auth.conf
[program:pdns_auth]
process_name=%(program_name)s
command=/usr/sbin/pdns_server
user=root
autostart=%(ENV_ENABLE_AUTH)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
